# local.d/cpu_performance.start
#!/bin/sh
echo 55000 | tee /sys/devices/virtual/thermal/thermal_zone0/trip_point_0_temp
